# etc/iso/crontab
@reboot root sleep 60 && /usr/local/bin/install.sh